# docker/images/pulse/supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; Entivia supervisord — manages all app processes inside the self-hosted image.
; Redis is started by docker-entrypoint.sh before supervisord takes over.
; nginx and Next.js are started here after the API is ready.
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/entivia

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── API — main REST server ───────────────────────────────────────────────────
[program:api]
command=uvicorn app.api.app:app
    --host 0.0.0.0
    --port 8000
    --workers %(ENV_API_WORKERS)s
    --loop uvloop
    --http httptools
directory=/app
user=entivia
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Worker — pipeline execution (pulls jobs from Redis queue) ────────────────
[program:worker]
command=python -m app.worker
directory=/app
user=entivia
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Agent — conversational AI + SSE streaming ────────────────────────────────
[program:agent]
command=uvicorn app.conversational.app:app
    --host 0.0.0.0
    --port 8001
    --workers %(ENV_AGENT_WORKERS)s
    --loop uvloop
    --http httptools
directory=/app
user=entivia
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Scheduler — pushes pipeline jobs to queue on cron schedule ───────────────
[program:scheduler]
command=python -m app.services.schedulers.run
directory=/app
user=entivia
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Next.js — serves the frontend (standalone build) ─────────────────────────
; Listens on 127.0.0.1:3000 — nginx proxies public traffic to it.
[program:nextjs]
command=node /frontend/server.js
environment=PORT="3000",HOSTNAME="127.0.0.1"
user=entivia
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── nginx — unified public entry-point on port 80 ─────────────────────────────
; Routes /api/* → FastAPI:8000, everything else → Next.js:3000.
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=5
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
